# longlifecoin.service
[Unit]
Description=LongLifeCoin Next.js Application
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/path/to/your/project
Environment=NODE_ENV=production
Environment=PORT=3000
ExecStart=/usr/bin/npm start
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target